# cookbooks/jenkins-jnlp/templates/jenkins.agent.service.erb
[Unit]
Description= jenkins JNLP agent service
After=network.target

[Service]
Type=simple
Restart=always
RestartSec=5
User=<%=node['jenkins_jnlp']['agent_owner']%>
WorkingDirectory=<%=node['jenkins_jnlp']['agent_home']%>
ExecStart=java -jar <%=node['jenkins_jnlp']['agent_home']%>/agent.jar -jnlpUrl http://<%=node['jenkins_jnlp']['agent_master']%>/computer/agent/jenkins-agent.jnlp -secret <%=@jenkins_key%> -workDir "<%=node['jenkins_jnlp']['var_path']%>"

[Install]
WantedBy=multi-user.target
