# cookbooks/jenkins-jnlp/templates/jenkins.agent.service.erb
[Unit]
Description= jenkins JNLP agent service
After=network.target

[Service]
Type=simple
Restart=always
RestartSec=5
User=root
WorkingDirectory=/home/vagrant
ExecStart=java -jar /home/vagrant/agent.jar -jnlpUrl http://jenkins-master:8080/computer/agent/jenkins-agent.jnlp -secret 742e122ffab5dec8810f5aa46ebbc292ce54638d4bc288fab2c26d1ef379fbb4 -workDir "/var/lib/jenkins"

[Install]
WantedBy=multi-user.target
